# stage/usr/lib/x86_64-linux-gnu/systemd/system/limine-manager-refresh.service
[Unit]
Description=Refresh Limine configuration generated by limine-manager
Documentation=man:limine-manager(8)

[Service]
Type=oneshot
ExecStart=/usr/bin/limine-manager --log-format=json refresh
